# commander.socket
[Socket]
ListenFIFO=/run/commander.sk

[Install]
WantedBy=sockets.target